# template/etc/systemd/system/dnsmasq.service.d/override.conf
[Service]
ExecStartPre=/opt/dnsmasq_generate_hosts
